# jp-post-agent/com.jexport.jppost.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jexport.jppost.sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__HOME__/Documents/jp-post-agent/sync.sh</string>
    </array>

    <!-- 30秒ごとに実行 -->
    <key>StartInterval</key>
    <integer>30</integer>

    <!-- ログイン時にも実行 -->
    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__HOME__/Documents/jp-post-agent/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Documents/jp-post-agent/launchd.err.log</string>
</dict>
</plist>
